# data/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:service-octane]
command=/usr/local/bin/php /var/www/artisan --env=$APP_ENV octane:start --server=swoole --port=8000
stdout_logfile=/var/log/service-octane.log
stderr_logfile=/var/log/service-octane.log
redirect_stderr=true
logfile=/dev/null
autorestart=true
autostart=true
priority=120
numprocs=1

[program:service-nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
stdout_logfile=/var/log/service-nginx.log
stderr_logfile=/var/log/service-nginx.log
logfile=/var/log/service-nginx.log
autorestart=false
autostart=true
priority=150
startsecs=0
numprocs=1
