# docker/supervisor/supervisor.d/server.conf

[program:app]
command=/usr/local/bin/server-start --host=0.0.0.0 --port=8000 --quiet
stdout_logfile=/var/log/%(program_name)_stdout.log
stderr_logfile=/var/log/%(program_name)_stderr.log
redirect_stderr=true
autorestart=true
stopasgroup=true
killasgroup=true
autostart=true
priority=120
numprocs=1
