# data/services/service-deploy.conf

[program:service-deploy]
command=/bin/sh /usr/local/bin/service-deploy
stdout_logfile=/var/log/service-deploy.log
stderr_logfile=/var/log/service-deploy.log
logfile=/dev/null
autorestart=true
autostart=true
priority=80
numprocs=1
