# data/services/service-schedule.conf

[program:service-schedule]
command=/usr/local/bin/php /var/www/artisan --env=$APP_ENV schedule:work
stdout_logfile=/var/log/service-schedule.log
stderr_logfile=/var/log/service-schedule.log
logfile=/dev/null
autorestart=true
autostart=true
priority=130
numprocs=1
