# data/services/service-fpm.conf

[program:service-fpm]
command=/usr/local/sbin/php-fpm -F -y /usr/local/etc/php-fpm.d/www.conf -c /usr/local/etc/php/php.ini
stdout_logfile=/var/log/service-fpm.log
stderr_logfile=/var/log/service-fpm.log
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
logfile=/dev/null
autorestart=true
autostart=true
priority=110
numprocs=1