# docker/services/service-watcher.conf

[program:watcher]
command=/usr/local/bin/watcher
stdout_logfile=/var/log/watcher_stdout.log
stderr_logfile=/var/log/watcher_stderr.log
autorestart=true
autostart=true
priority=150
numprocs=1
