# docker/supervisor/conf.d/app.conf

[program:app]
command=/usr/local/bin/artisand --host=0.0.0.0 --port=3300 --admin-port=2019 --quiet
stdout_logfile=/var/log/%(program_name)-stdout.log
stderr_logfile=/var/log/%(program_name)-stderr.log
stopsignal=SIGTERM
directory=/var/www
autorestart=true
stopasgroup=true
killasgroup=true
autostart=true
priority=120
numprocs=1
